# rica.service
[Unit]
Description=Rica Production Stack
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/path/to/rica
ExecStart=/bin/bash /path/to/rica/deploy-prod.sh start
ExecStop=/bin/bash /path/to/rica/deploy-prod.sh stop
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
